# image-upload/services/agrodrone-image-upload.path
[Unit]
Description=Watch for offload trigger post-mission
After=multi-user.target

[Path]
PathExists=/tmp/offload_requested

[Install]
WantedBy=multi-user.target
